# net/ntpsec/files/ntpdate.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2008 OpenWrt.org

START=60

STEP_SERVERS="$(uci get system.ntp.server)"
TIMEOUT="2" # in seconds

start() {
	for s in $STEP_SERVERS ; do
		/usr/bin/ntpdate -s -b -u -t "$TIMEOUT" "$s" && break
	done
}
